# docker-files/etc/supervisor/conf.d/nginx-epoll.conf
[program:nginx-epoll]
command = /usr/sbin/nginx -c /etc/nginx/nginx-epoll.conf
autostart = false
autorestart = unexpected
exitcodes = 0
