# docker-files/etc/supervisor/conf.d/uwsgi.conf
[program:uwsgi]
user = www-data
command = uwsgi --emperor --ini /etc/uwsgi-emperor/emperor.ini
process_name = %(program_name)s_%(process_num)02d
autostart = true
autorestart = unexpected
