# docker-files/etc/supervisor/conf.d/uwsgi.conf
[program:uwsgi]
user = www-data
command = /app/venv3.11/bin/launch-supervisor-uwsgi
process_name = %(program_name)s_%(process_num)02d
autostart = true
autorestart = unexpected
