# docker-files/etc/supervisor/conf.d/uwsgi.conf
[program:uwsgi]
command = /usr/bin/uwsgi --emperor --ini /etc/uwsgi-emperor/emperor.ini
user = www-data
process_name = %(program_name)s_%(process_num)02d
autostart = true
autorestart = unexpected
